# systemd/bot.service
[Service]
User=root
WorkingDirectory=/home/<root_directory>
ExecStart=/usr/bin/python3 main.py
Restart=always

[Install]
WantedBy=multi-user.target
